# 2-Server/sv/foundationdb/finish
#!/bin/sh

[ -r conf ] && . ./conf

rm -f "$PIDFILE"
