# 2-Server/sv/foundationdb/run
#!/bin/sh

[ -r conf ] && . ./conf

exec 2>&1
exec "$FDBMONITOR" --conffile "$CONFFILE" --lockfile "$PIDFILE"
